# rtmp_printer_streamer/rootfs/etc/services.d/streamer/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

log() { echo "$(date -Iseconds) [streamer] $*"; }

STREAM_KEY=$(bashio::config 'stream_key')
PRINTER_URL=$(bashio::config 'printer_url')
CHECK_INTERVAL=$(bashio::config 'check_interval')
STATUS_URL="$(echo "${PRINTER_URL}" | sed -E 's#^(https?://[^/]+).*#\1#')"

RTMP_BASE="rtmp://rtmp.deployor.dev:1935/live/"
RTMP_PUBLISH_URL="${RTMP_BASE}${STREAM_KEY}"

STREAMER_PID=""

is_streamer_running() {
  [[ -n "${STREAMER_PID:-}" ]] && kill -0 "${STREAMER_PID}" 2>/dev/null
}

start_stream() {
  log "Starting Python streamer: ${PRINTER_URL} -> ${RTMP_PUBLISH_URL}"
  python3 /usr/local/bin/stream.py "${PRINTER_URL}" "${RTMP_PUBLISH_URL}" &
  STREAMER_PID=$!
  log "Python streamer started with PID ${STREAMER_PID}"
}

stop_stream() {
  if is_streamer_running; then
    log "Stopping streamer PID ${STREAMER_PID}"
    kill "${STREAMER_PID}" 2>/dev/null || true
    wait "${STREAMER_PID}" 2>/dev/null || true
  fi
  STREAMER_PID=""
}

check_online() {
  # Consider online if a quick GET returns HTTP 200 on the base URL (host:port)
  curl -s --fail -o /dev/null --max-time 5 "${STATUS_URL}"
}

cleanup() {
  log "Termination signal received, cleaning up..."
  stop_stream
  exit 0
}

trap cleanup SIGTERM SIGINT

log "Starting monitor loop (interval: ${CHECK_INTERVAL}s)"

while true; do
  if check_online; then
    if ! is_streamer_running; then
      start_stream || { log "Error: failed to start streamer"; STREAMER_PID=""; }
    fi
  else
    if is_streamer_running; then
      log "Printer offline, stopping stream"
      stop_stream
    else
      log "Printer offline, waiting"
    fi
  fi

  # streamer died unexpectedly -> reset state; will restart next tick if online
  if [[ -n "${STREAMER_PID:-}" ]] && ! is_streamer_running; then
    log "streamer exited unexpectedly"
    STREAMER_PID=""
  fi

  # Monitor the stream
  while true; do
    if ! kill -0 $STREAMER_PID 2>/dev/null; then
      log "Stream process died, restarting..."
      break
    fi
    sleep 5
  done

  sleep "${CHECK_INTERVAL:-10}"

done